# deploy/systemd/freguesia-local-api.service
[Unit]
Description=Freguesia Local API
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/Freguesia
EnvironmentFile=/root/Freguesia/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/node /root/Freguesia/server/local-api.mjs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
